# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D -e

[program:fcgi]
command=/usr/bin/fcgiwrap -s unix:/var/run/fcgiwrap.sock

[program:nginx]
command=/usr/sbin/nginx
